# docs/NexGenHost — Go Live - infrastructure files/nexhost-worker.service
[Unit]
Description=NexGenHost Go Worker Engine
Documentation=https://github.com/nexgentech/nexgenhost
After=network.target docker.service
Requires=docker.service

[Service]
# ── Identity ──────────────────────────────────────────────────────────────────
Type=simple
User=nexhost
Group=docker
WorkingDirectory=/opt/nexhost-worker

# ── Binary ────────────────────────────────────────────────────────────────────
ExecStart=/opt/nexhost-worker/bin/worker

# ── Environment ───────────────────────────────────────────────────
Environment="WORKER_REDIS_URL=rediss://default:gQAAAAAAAlYGAAIgcDI1MjRlYjM5Y2ZlYTk0YTk3YTA1NWQ4MjE5MGE0ODE4Mg@summary-wildcat-153094.upstash.io:6379"
Environment="WORKER_DEPLOY_QUEUE=nexgenhost-deployments"
Environment="WORKER_API_CALLBACK_URL=https://cloud-platform-5vf4.onrender.com/internal/deploy/callback"
Environment="INTERNAL_SECRET=6da37ed6006c95bb1eecffdbeba4992c279e8181645d522145bf07333e58b832"
Environment="DOCKER_HOST=unix:///var/run/docker.sock"
Environment="CONTAINER_NETWORK=nexhost-bridge"
Environment="WORKER_WORK_DIR=/opt/nexhost-worker/tmp"
Environment="WORKER_ID=worker-1"
Environment="LOG_LEVEL=info"

# ── Restart policy ────────────────────────────────────────────────────────────
Restart=always
RestartSec=5s
StartLimitInterval=60s
StartLimitBurst=3

# ── Logging ───────────────────────────────────────────────────────────────────
StandardOutput=append:/opt/nexhost-worker/logs/worker.log
StandardError=append:/opt/nexhost-worker/logs/worker-error.log

# ── Resource limits ───────────────────────────────────────────────────
# Oracle Free Tier: VM.Standard.E2.1.Micro — 1 OCPU, 1GB RAM (+4GB swap)
# Worker stays lean; most resources reserved for Docker build containers
MemoryMax=256M
CPUQuota=80%

# ── Security hardening ────────────────────────────────────────────────────────
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/opt/nexhost-worker /var/run/docker.sock /tmp

[Install]
WantedBy=multi-user.target
